# muteads.service.template
[Unit]
Description="mute Spotify ads"

[Service]
ExecStart=spotify-scripts/muteads.sh

[Install]
WantedBy=multi-user.target
